# chappie-training.service
[Unit]
Description=CHAPPiE Autonomous Training Service
After=network.target

[Service]
# User unter dem es laufen soll (Anpassen falls nötig!)
User=bbecker
Group=bbecker

# Pfad zum Projekt
WorkingDirectory=/home/bbecker/CHAPPiE

# Environment Variablen laden
Environment="PATH=/home/bbecker/CHAPPiE/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
# Falls du Cerebras Token direkt hier setzen willst:
# Environment="CEREBRAS_API_KEY=dein-key"

# Start-Befehl (Direkt den Loop starten, nutzt gespeicherte Config)
ExecStart=/home/bbecker/CHAPPiE/venv/bin/python3 -m Chappies_Trainingspartner.training_loop

# Restart bei Crash (Autonomie ist Key!)
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target