# chappie-vllm.service
[Unit]
Description=CHAPPiE Steering OpenAI API Server
After=network.target

[Service]
User=bbecker
Group=bbecker
WorkingDirectory=/home/bbecker/CHAPPiE

Environment="PATH=/home/bbecker/CHAPPiE/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONIOENCODING=utf-8"
Environment="CUDA_VISIBLE_DEVICES=0"
Environment="HF_HOME=/home/bbecker/.cache/huggingface"
Environment="CHAPPIE_STEERING_MODEL=Qwen/Qwen3-4B-Instruct-2507"

ExecStart=/home/bbecker/CHAPPiE/venv/bin/python3 -m brain.steering_api_server --host 127.0.0.1 --port 8000 --model Qwen/Qwen3-4B-Instruct-2507

Restart=always
RestartSec=10
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target