# deploy/chappie-web.service
[Unit]
Description=CHAPPiE App API
After=network.target chappie-vllm.service
Wants=chappie-vllm.service

[Service]
User=bbecker
Group=bbecker
WorkingDirectory=/home/bbecker/CHAPPiE
Environment="PATH=/home/bbecker/CHAPPiE/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONIOENCODING=utf-8"
Environment="CUDA_VISIBLE_DEVICES=-1"
Environment="CHROMA_SQLITE_JOURNAL_MODE=WAL"
Environment="TOKENIZERS_PARALLELISM=false"
Environment="OMP_NUM_THREADS=1"
ExecStart=/home/bbecker/CHAPPiE/venv/bin/python -m uvicorn api.main:app --host 0.0.0.0 --port 8010
Restart=always
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
